# package/luci/applications/luci-app-netdata/root/etc/uci-defaults/40_luci-app-netdata
#!/bin/sh

[ -d /usr/share/netdata/zh_Hant ] && cp -f /usr/share/netdata/zh_Hant/* /usr/share/netdata/web/

rm -rf /tmp/luci-modulecache /tmp/luci-indexcache*

exit 0
